# hadoop-dev/profile.d/hadoop.sh
#!/usr/bin/env bash

export HADOOP_HOME=/hadoop
