# devops-python-tools-fedora/profile.d/java.sh
#!/usr/bin/env bash

export JAVA_HOME=/usr
