# apache-drill/supervisord.d/drill.ini
[program:drill]
#command=/apache-drill/bin/drillbit.sh start && sleep 2 && tail -f /apache-drill/log/*
command=/apache-drill/bin/drillbit.sh run
