# centos-dev/profile.d/jython.sh
export JYTHON_HOME=/opt/jython
export PATH=$PATH:$JYTHON_HOME/bin
