# ubuntu-dev/profile.d/gradle.sh
#!/usr/bin/env bash

export GRADLE_HOME=/opt/gradle
export PATH=$PATH:$GRADLE_HOME/bin
