# apache-drill/drill.ini
[program:drill]
command=/apache-drill/bin/drillbit.sh start && cat
