# alpine-dev/profile.d/groovy.sh
export GROOVY_HOME=/opt/groovy
export PATH=$PATH:$GROOVY_HOME/bin
