# alpine-dev/profile.d/sbt.sh
export SBT_HOME=/opt/sbt
export PATH=$PATH:$SBT_HOME/bin
